# data/Dockerfiles/syncengine/supervisord.conf
[supervisord]
nodaemon=true

[program:sync_engine_inbox]
priority=1
command=/etc/sync-engine/inbox-start/run
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
startsecs=0

[program:sync_engine_api]
priority=1
command=/etc/sync-engine/inbox-api/run
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
startsecs=0
